# src/profile.d/base
# Add ~/.bin to the path
# This is where we will keep custom executables
export PATH=$HOME/.bin:/usr/local/bin:/usr/local/sbin:$PATH

# Set Unicode stuff
export LANG=en_US.UTF-8
export LC_ALL=en_US.UTF-8
export LC_CTYPE=en_US.UTF-8

# Check the window size after each command and, if necessary,
# update the values of LINES and COLUMNS.
shopt -s checkwinsize

