# src/profile.d/clojure
export PATH=~/.cljr/bin:$PATH

# vim: ft=sh
