# src/profile.d/dev
# set default editor
export EDITOR=gvim_open

# vi: ft=sh
