# src/profile.d/dev.erb
# set default editor
<%- if RUBY_PLATFORM =~ /darwin/ -%>
export EDITOR=mate
<%- else -%>
export EDITOR=vi
<%- end -%>

if which brew >/dev/null 2>&1 && test -f $(brew --prefix)/etc/autojump; then \
  . $(brew --prefix)/etc/autojump
fi

