# src/profile.d/etc
function timestamp {
  date +"%Y%m%d%H%M%S"
}

# vi: ft=sh
