# src/profile.d/handy_shortcuts
# enable color support of ls
alias ls='ls --color=auto -p'

# some more ls aliases
#alias ll='ls -lT'
alias ll='ls -l'
alias llt='ls -lTt'
alias lv='ls -1'
#alias lli="ll -X | tac | head -n -1 | tac | isort -i '/$' ' ([^ ]+)$'"  # kind of naive but decent.
alias lsh="ls -a"
alias llh="ll -a"

# A more handy rgrep
smartgrep() {
   grep -r -E -n "$@" . | grep -E -v '\.svn|~:[[:digit:]]+:'
}

# Make it so ANSI colors will still show up in less
alias less='less -r'
