# src/profile.d/mysql
# Add MySQL stuff to path
export PATH=/usr/local/mysql/bin:$PATH

# vi: ft=sh
