# src/profile.d/osx
# http://refactormycode.com/codes/63-open-a-new-tab-in-current-directory-from-iterm
# https://bitbucket.org/vpetro/applescript/raw/7f3ac525ffa0/term
nt()
{
  osascript -e "
  tell application \"iTerm\"
    activate
    tell the current terminal
      launch session \"Default Session\"
      tell the last session
        write text \"cd $(pwd)\"
      end tell
    end tell
  end tell"
}
