# src/profile.d/ruby
# Load Rubygems every time Ruby loads
export RUBYOPT=rubygems

ri() {
   ri -T -f ansi "$@" | less -r
}
alias gi="sudo gem install $@ --no-ri --no-rdoc"
