# src/profile.d/vim
alias e='gvim_open_file'
alias d='gvim_open_dir'

# vim: ft=sh
