# profile.d/code_swarm
# add code_swarm to path
export PATH=/opt/code_swarm/bin:$PATH
