# profile.d/dev
# set default editor
export EDITOR=mate
#export EDITOR=vi

#alias gotoproj="cd ~/code/work/premiere/proj"
#alias goto="source ~/code/common/scripts/bash/find_proj.sh"
#export CDPATH=$HOME/work/premiere/proj:$HOME/personal/proj/web:$HOME/code/github:$HOME/code/aux

alias ey-amtees="EYRC=~/code/work/premiere/proj/amtees/.eyrc ey"
alias ey-pcoll="EYRC=~/code/work/premiere/proj/collectibles/.eyrc ey"

if [ -f `brew --prefix`/etc/autojump ]; then \
  . `brew --prefix`/etc/autojump
fi

