# profile.d/ec2
# load amazon keys
if [[ -f "$HOME/.amazon_keys" ]]; then
  . "$HOME/.amazon_keys";
fi
# ec2-utils stuff
export PATH=/opt/ec2-tools/bin:$PATH
export EC2_HOME=/opt/ec2-tools
export EC2_CERT=~/etc/cert-RDVVMUU7Z4HTH4RORRG4JNCAJJLG42S5.pem
export EC2_PRIVATE_KEY=~/etc/pk-RDVVMUU7Z4HTH4RORRG4JNCAJJLG42S5.pem

