# profile.d/long_running_commands
# from <http://www.macosxhints.com/article.php?story=20071009124425468>

export PREEXEC_EXCLUDE_LIST="mc less bash man nano say script/console script/server irb svndiff svn spork screen"
export PREEXEC_MAX=10

. ~/.preexec.bash

# called before each command and starts stopwatch
function preexec () {
	export PREEXEC_CMD=`echo $1 | awk '{ print $1; }'`
	export PREEXEC_TIME=$(date +'%s')
}

# called after each command, stops stopwatch
# and notifies if time elpsed exceeds threshold
function precmd () {
	stop=$(date +'%s')
	start=${PREEXEC_TIME:-$stop}
	let elapsed=$stop-$start
	max=${PREEXEC_MAX:-10}
	
	for i in ${PREEXEC_EXCLUDE_LIST:-}; do
	  if [ "x$i" = "x$PREEXEC_CMD" ]; then
	    max=999999;
	    break;
	  fi
	done
	
	## CURRENTLY THIS FIRES BY MISTAKE IF YOU CTRL-C ON AN EMPTY LINE
	#if [ $elapsed -gt $max ]; then
	#	growlnotify -n "Terminal" -m "is now finished (time: $elapsed secs)" "${PREEXEC_CMD:-a terminal command}"
        #        #say "${PREEXEC_CMD:-your terminal command} is now finished: (time: $elapsed seconds)"
        #        afplay /System/Library/Sounds/Blow.aiff
	#fi
}

preexec_install
