# profile.d/rvm
# rvm installer added line:
if [ -s ~/.rvm/scripts/rvm ] ; then source ~/.rvm/scripts/rvm ; fi

# Thanks for the awesome idea batasrki function
function gemdir {
  if [ -z "$1" ] ; then
    echo "gemdir expects a parameter, which should be a valid rvm Ruby selector"
  else
    rvm use $1
    cd `rvm gemdir`
    pwd
  fi
}

export PS1_RVM="\[\033[01;33m\]\$(~/.rvm/bin/rvm-prompt)\[\033[00m\]"
export PS1_START="$PS1_START $PS1_RVM"
export PS1="$PS1_START$PS1_END"
