# profile.d/svn
function svnstatus() {
  svn status "$@" | grep -v '.[ ]\+\(log\|tmp\)' | sort -k 1,1
}
alias svnst="svnstatus"
# Color diff
# <http://www.pjhyett.com/posts/203-colored-svn-diff>
#function svndiff() {
#  svn diff --diff-cmd=colordiff "$@" | less
#}
svnstatusinline() {
  svnstatus | grep "^$1" | cut -c 8- | argize
}
svnci() {
  EDITOR=nano svn commit "$@"
}

rmsvnfiles() {
  find $1 -name .svn -exec rm -rfv {} \;
}

